# actividad4/actividad.service
[Unit]
Description=Servicio de Saludo

[Service]
ExecStart=/usr/local/bin/saludo.sh
Restart=always

[Install]
WantedBy=multi-user.target